# deploy/cellyn-store.service
[Unit]
Description=Cellyn Store (Node.js + SQLite)
After=network.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=NODE_ENV=production
Environment=PORT=__PORT__
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Hardening ringan
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
